# www-apps/lidarr/files/lidarr.init
#!/sbin/openrc-run

pidfile="/run/lidarr.pid"
command="/opt/lidarr/Radarr"
command_user="lidarr:lidarr"
command_background=true

